# files/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:start_tinydns]
command=/start_tinydns.bash
startretries=0
startsecs=0
autorestart=false

